# systemd/vaultctl-renew.timer
[Unit]
Description=Vault Token Auto Renew Timer

[Timer]
# 매시간 실행
OnCalendar=hourly

# 시스템 부팅 후 5분 뒤 첫 실행
OnBootSec=5min

# 놓친 실행 복구
Persistent=true

# 랜덤 지연 (여러 서버에서 동시 실행 방지)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
